# systemd/miner.timer
[Unit]
Description=Start miner

[Timer]
OnBootSec=1sec
Persistent=true

[Install]
WantedBy=timers.target